# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=apache2ctl -DFOREGROUND
killasgroup=true
stopasgroup=true
autorestart=true